# mariadb/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command:/usr/sbin/mysqld --user=root

